# azure.service
[Service]
ExecStart=/bin/python3 %h/azure.py
